# deploy/macos/io.macvendor.offhost-backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>io.macvendor.offhost-backup</string>
  <key>ProgramArguments</key>
  <array><string>__LIBEXEC__/macvendor-offhost-backup</string></array>
  <key>StartCalendarInterval</key>
  <dict><key>Hour</key><integer>4</integer><key>Minute</key><integer>20</integer></dict>
  <key>ProcessType</key><string>Background</string>
  <key>StandardOutPath</key><string>__LOGS__/offhost-backup.log</string>
  <key>StandardErrorPath</key><string>__LOGS__/offhost-backup.err.log</string>
</dict>
</plist>
