# deploy/macvendor-backup.timer
[Unit]
Description=Daily macvendor PostgreSQL backup

[Timer]
OnCalendar=*-*-* 02:20:00 UTC
RandomizedDelaySec=10m
Persistent=true
Unit=macvendor-backup.service

[Install]
WantedBy=timers.target
