# deploy/macvendor-source-update.timer
[Unit]
Description=Daily atomic macvendor source update

[Timer]
OnCalendar=*-*-* 02:10:00 UTC
RandomizedDelaySec=15m
Persistent=true
Unit=macvendor-source-update.service

[Install]
WantedBy=timers.target
