# deploy/macvendor-traffic-report.timer
[Unit]
Description=Refresh macvendor traffic report every 15 minutes

[Timer]
OnBootSec=5m
OnUnitActiveSec=15m
RandomizedDelaySec=2m
Persistent=true
Unit=macvendor-traffic-report.service

[Install]
WantedBy=timers.target
